# deploy/exclaw.service
# ExClaw systemd service — bare-metal Elixir release on DGX Spark.
#
# Install:
#   sudo cp deploy/exclaw.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable exclaw
#
# Usage:
#   sudo systemctl start exclaw
#   sudo systemctl stop exclaw
#   journalctl -u exclaw -f
#   /opt/exclaw/current/bin/exclaw remote   # attach IEx

[Unit]
Description=ExClaw Personal AI Assistant
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=wido
Group=wido

WorkingDirectory=/opt/exclaw
EnvironmentFile=/opt/exclaw/.env
Environment=MIX_ENV=prod
Environment=LANG=en_US.UTF-8

# PHX_SERVER=true is set by bin/server overlay.
ExecStart=/opt/exclaw/current/bin/server

# BEAM handles SIGTERM gracefully — give it time to drain.
KillMode=mixed
TimeoutStopSec=30

Restart=on-failure
RestartSec=5

# Docker socket access for Container.Manager sandbox spawning.
SupplementaryGroups=docker

LimitNOFILE=65536

# Hardening — minimize attack surface.
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectHostname=true
ProtectClock=true
RestrictSUIDSGID=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=exclaw

[Install]
WantedBy=multi-user.target
